# debian/leicht.postinst
#!/bin/sh

systemctl daemon-reload
systemctl enable leicht.target
systemctl restart leicht.target

#DEBHELPER#
