# src/init/leicht.target
[Unit]
Description=Leicht - Telegram bot framework
After=multi-user.target

[Install]
WantedBy=multi-user.target
